# systemd/user/waybar.service
[Unit]
Description=Waybar

[Service]
ExecStart=/run/current-system/sw/bin/waybar
Restart=on-failure
RestartSec=1
